# etc/systemd/system/willie.service.d/local.conf
[Service]
# Override config file used
ExecStart=
ExecStart=/data/applications/grobbebol/willie/willie.py -c /data/applications/grobbebol/grobbebol.cfg
Restart=on-failure
