# etc/supervisor/conf.d/exit2020.conf
[program:exit2020]
directory=/data/www/fizzgig.nl/applications/exit2020/src/
command=/data/www/fizzgig.nl/applications/exit2020/venv/bin/python -m exit 2020
stopsignal=INT
killasgroup=true
user=bpeschier
redirect_stderr=true
stdout_logfile=/data/www/fizzgig.nl/applications/exit2020/logs/exit2020.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
