# etc/systemd/system/getty@.service.d/local.conf
[Unit]
# Override check for /dev/tty0, which a container does not have
ConditionPathExists=
